# crontab
# “At 02:00 on every 7th day-of-month from 1 through 31”
0 2 1,7,14,21,28 * * /scripts/run_certbot.sh >> /var/log/letsencrypt/cron.log 2>&1
