# crontab
0 2 1/7 * * /scripts/run_certbot.sh >> /var/log/letsencrypt.log 2>&1
