# crontab
* * * * * /scripts/run_certbot.sh >> /var/log/letsencrypt.log 2>&1
