# crontab
# “At 02:00 on every day-of-month”
0 2 * * * /scripts/run_certbot.sh >> /var/log/letsencrypt/cron.log 2>&1
